# systemd/das-backup-full.service
[Unit]
Description=DAS Full Backup - Recreate boot subvolumes from latest snapshots
Documentation=file:///hddRaid1/ClaudeCodeProjects/DAS-Backup-Manager/docs/OFFLINE-BACKUP-PLAN.md
After=local-fs.target
ConditionPathExistsGlob=/dev/disk/by-id/usb-TerraMas_TDAS*

[Service]
Type=oneshot
ExecStart=/hddRaid1/ClaudeCodeProjects/DAS-Backup-Manager/scripts/backup-run.sh --full
StandardOutput=journal
StandardError=journal
Nice=19
IOSchedulingClass=idle
TimeoutStartSec=21600

[Install]
WantedBy=multi-user.target
